# etc/init/cm160.conf
# CM160 daemon is part of eagle-owl, which is meant to interface the Owl cm160 +USB energy monitor

description	"CM160 server"
author          "Frederic Bourqui <fbourqui@yahoo.com>"

start on local-filesystems or runlevel [2345]
stop on runlevel [!2345]

expect fork
respawn
respawn limit 10 5

pre-start script
    test -x /opt/eagleowl/cm160 || { stop; exit 0; }
end script

script
    cd /opt/eagleowl/
    exec /opt/eagleowl/cm160 -d
end script
