# etc/init/cm160.conf
# CM160 daemon is part of eagle-owl, which is meant to interface the Owl cm160 +USB energy monitor

description	"CM160 server"
author          "Frederic Bourqui <fbourqui@yahoo.com>"

start on local-filesystems or runlevel [2345]
stop on runlevel [!2345]

expect fork
respawn
respawn limit 10 5

pre-start script
    . /etc/eagleowl.conf
    test -x $install_path/cm160 || { stop; exit 0; }
end script

script
    . /etc/eagleowl.conf
    cd $install_path
    exec $install_path/cm160 -d
end script
